# tests/vault/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[unix_http_server]
file = /tmp/supervisor.sock

[program:spire-agent]
command = /opt/spire/bin/spire-agent run -config /opt/spire/conf/agent/agent.conf

[program:vault]
command = vault server -config /vault/config/vault-local.hcl

